# tools/fio.service
[Unit]

Description=flexible I/O tester server
After=network.target

[Service]

Type=simple
PIDFile=/run/fio.pid
ExecStart=/usr/bin/fio --server
